# crontab
0 0 * * * logrotate.sh
* * * * * callback_sender.php
